# php-fpm-base/supervisord.d/crond.conf
[program:crond]
process_name=%(program_name)s
command=cron -f -L 15
autostart=true
autorestart=true
