# php-fpm-base/supervisord.d/php-fpm.conf
[program:php-fpm]
process_name=%(program_name)s
command=php-fpm
autostart=true
autorestart=true
